# etc/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
history_file=~/.sc_history

[program:perlpad]
directory=/home/homepage/PerlPad
command=/home/homepage/PerlPad/etc/env.sh perl -Mlib=extlib/lib/perl5 extlib/bin/plackup -s Starman --workers 30 -E deployment --preload-app app.psgi
user=homepage
autostart=true
autorestart=true
startsecs=1
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/perlpad.log
stderr_logfile=/var/log/perlpad.log
environment=HOME="/home/homepage",USER="homepage",PLACK_ENV="<MODE>"

;[include]
;files=/etc/supervisord.d/*.ini

