# config/supervisord/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
history_file=~/.sc_history

[program:perlpad]
directory=/home/homepage/GrowthPerl
user=homepage
autostart=true
autorestart=true
startsecs=1
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/GrowthPerl/perlpad.log
stderr_logfile=/var/log/GrowthPerl/perlpad.log
environment=HOME="/home/homepage",USER="homepage"
command=perl -Mlib=extlib/lib/perl5 extlib/bin/plackup -s Starman -E deployment --preload-app --disable-keepalive --workers 10 app.psgi

;command=/home/homepage/GrowthPerl/etc/env.sh perl -Mlib=extlib/lib/perl5 extlib/bin/plackup -s Starman --workers 30 -E deployment --preload-app app.psgi
;command=perl -Mlib=extlib/lib/perl5 extlib/bin/plackup -s Starman --workers 30 -E deployment --preload-app app.psgi

;[include]
;files=/etc/supervisord.d/*.ini
